# DEBIAN/postinst
pip install bluesky zmq systemd
# reload
systemctl daemon-reload
# enable the service
systemctl enable bluesky-zmq.service
systemctl start bluesky-zmq.service
# to stop:
systemctl stop bluesky-zmq.service
